# mediacenter/data/etc/systemd/system/devmon.service
[Unit]
Description=devmon

[Service]
User=root
Group=root
Type=simple
ExecStart=/usr/bin/devmon

[Install]
WantedBy=multi-user.target
